# overlay/etc/profile.d/qtwebengine.sh
export QTWEBENGINE_CHROMIUM_FLAGS="--disable-es3-gl-context --ignore-gpu-blacklist"
